# automation/com.tommy.forgejo-webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tommy.forgejo-webhook</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/nqcdan/dev/wiki/automation && source .venv/bin/activate && exec uvicorn forgejo_webhook_service:app --host 0.0.0.0 --port 9009</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/nqcdan/dev/wiki/automation/forgejo-webhook.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nqcdan/dev/wiki/automation/forgejo-webhook.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>WEBHOOK_SECRET</key>
    <string>nqcdan1908</string>
    <key>GIT_REMOTE</key>
    <string>origin</string>
    <key>GIT_BRANCH</key>
    <string>main</string>
  </dict>
</dict>
</plist>
